# compusona-boot.service.example
[Unit]
Description=compusona greeting on boot
After=network-online.target
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/compusona boot

[Install]
WantedBy=multi-user.target
